# examples/systemd/mxcp-audit-cleanup.timer
[Unit]
Description=Run MXCP Audit Cleanup daily at 2 AM
Requires=mxcp-audit-cleanup.service

[Timer]
# Run daily at 2:00 AM
OnCalendar=daily
AccuracySec=1h
Persistent=true

# Randomize start time by up to 30 minutes to avoid thundering herd
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
